# metrics_role/templates/metrics_exporter.service.j2
[Unit]
Description=Metrics Exporter Service
After=network.target

[Service]
ExecStart=/usr/local/bin/metrics_exporter.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target
